# ours/chatnest/full-stack/ethan-chat.service.example
[Unit]
Description=Ethan ChatNest full-stack
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/ethan-bot/ours/chatnest/full-stack
Environment=HOST=127.0.0.1
Environment=PORT=8787
ExecStart=/srv/ethan-bot/ours/chatnest/full-stack/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8787
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
